# src/vendor/Anakeen/anakeen.cron
# Every 5 minutes
*/5 * * * * [ -f "$CONTEXT_ROOT"/ank.php ] && [ ! -f "$CONTEXT_ROOT"/maintenance.lock ] && cd "$CONTEXT_ROOT" && ./ank.php --script=processExecute > /dev/null